# examples/systemd-cove-cinema.service.example
# systemd unit — Linux / VPS 上跑 cove-cinema。
#
# 装：
#   1. 改下面三处 REPLACE_*
#   2. sudo cp 这个文件到 /etc/systemd/system/cove-cinema.service
#   3. sudo systemctl daemon-reload
#   4. sudo systemctl enable --now cove-cinema
#   5. 看日志：journalctl -u cove-cinema -f
#   6. 重启：sudo systemctl restart cove-cinema
#
# 配 nginx 反代示例（前端不直连 8770，走 https）：
#   location /cinema/ {
#       proxy_pass http://127.0.0.1:8770;
#       proxy_http_version 1.1;
#       proxy_set_header Host $host;
#       proxy_set_header X-Real-IP $remote_addr;
#       # 视频流大文件别 buffer，要快速吐字节
#       proxy_buffering off;
#       proxy_request_buffering off;
#       # 不要 timeout 视频流：read timeout 拉大
#       proxy_read_timeout 3600s;
#   }

[Unit]
Description=cove-cinema video streaming backend
After=network.target

[Service]
Type=simple
User=REPLACE_WITH_USERNAME
WorkingDirectory=REPLACE_WITH_FULL_PATH/cove-cinema
ExecStart=/usr/bin/python3 REPLACE_WITH_FULL_PATH/cove-cinema/server.py --port 8770 --root REPLACE_WITH_VIDEOS_PATH
Restart=on-failure
RestartSec=3

# 想直接对外暴露（不走 nginx）就加 --bind 0.0.0.0：
# ExecStart=/usr/bin/python3 .../server.py --port 8770 --bind 0.0.0.0 --root /data/films

[Install]
WantedBy=multi-user.target
